# cookbooks/hive/templates/default/hive-server.erb
#!/bin/bash
# init script for Hive Thrift Interface.
#
# chkconfig: 2345 90 10
# description: Hive Thrift Interface

# Source function library.
. /etc/rc.d/init.d/functions

# Paths to configuration, binaries, etc
HIVE_BIN="<%= node[:hive][:home_dir] %>/bin/hive"
HIVE_ARGS="--service hiveserver"
HIVE_LOG="<%= node[:hive][:log_dir] %>/hive-thrift.log"
HIVE_USER="<%= node[:hive][:user] %>"
ANT_LIB=/usr/share/java

if [ ! -f $HIVE_BIN ]; then
  echo "File not found: $HIVE_BIN"
  exit 1
fi

# pid file for /sbin/runuser
pidfile=${PIDFILE-<%= node[:hive][:pid_dir] %>/hive-thrift.pid}
RETVAL=0

start() {
  # check to see if hive is already running by looking at the pid file and grepping
  # the process table.
  if [ -f $pidfile ] && checkpid `cat $pidfile`; then
    echo "hive-thrift is already running"
    exit 0
  fi
  echo -n $"Starting $prog: "
  nohup sudo -u $HIVE_USER $HIVE_BIN $HIVE_ARGS >> $HIVE_LOG 2>&1 &
  java_pid=$!
  echo $java_pid > $pidfile
  ps $java_pid > /dev/null
  return $?
}

stop() {
  # check if the process is already stopped by seeing if the pid file exists.
  if [ ! -f $pidfile ]; then
    echo "hive-thrift is already stopped"
    exit 0
  fi
  echo -n $"Stopping $prog: "
  if kill `cat $pidfile`; then
    RETVAL=0
    echo_success
  else
    RETVAL=1
    echo_failure
  fi
  echo
  [ $RETVAL = 0 ] && rm -f ${pidfile}
}

status_fn() {
  if [ -f $pidfile ] && checkpid `cat $pidfile`; then
    echo "hive-thrift is running"
    exit 0
  else
    echo "hive-thrift is stopped"
    exit 1
  fi
}

case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  status)
    status_fn
    ;;
  restart)
    stop
    start
    ;;
  *)
    echo $"Usage: $prog {start|stop|restart|status}"
    RETVAL=3
esac

exit $RETVAL
